# laravel-horizon/supervisord.d/auth-notification.conf
[program:auth-notification]

process_name=%(program_name)s_%(process_num)02d
command=php /var/www/auth/artisan rabbitevents:listen Notification.#
autostart=true
autorestart=true
redirect_stderr=true