# laravel-horizon/supervisord.d/listen-order.conf
[program:listen-order]

process_name=%(program_name)s_%(process_num)02d
command=php /var/www/prod/_work/tomoni-product/tomoni-product/artisan rabbitevents:listen Order.#
autostart=true
autorestart=true
redirect_stderr=true
