# laravel-horizon/supervisord.d/listen-product.conf
[program:listen-product]

process_name=%(program_name)s_%(process_num)02d
command=php /var/www/prod/_work/tomoni-auth/tomoni-auth/artisan rabbitevents:listen Product.#
autostart=true
autorestart=true
redirect_stderr=true
